# bspwm/bin/unimatrix-screensaver.service
[Unit]
Description=Screensaver

[Service]
User=venerable_white
Type=simple
ExecStart=/bin/bash /home/venerble_white/.config/bspwm/bin/unimatrix-ss

[Install]
WantedBy=multi-user.target